# s3sync/scripts/s3sync.timer
[Unit]
Description=Periodically run s3sync

[Timer]
OnUnitActiveSec=
OnCalendar=*:*:0/30
Persistent=true

[Install]
WantedBy=timers.target
